# modules/e2/templates/etc/rc.local.erb
#!/bin/sh -e

<% @init.each do |command| %>
<%= command %>
<% end %>

exit 0
