# modules/init/templates/etc/rc.local.erb
#!/bin/sh -e

<% @init.each do |command| %>
<%= command %>
<% end %>

if [ -f /etc/local.d/startup.sh ] then
	/etc/local.d/startup.sh
fi

exit 0
